# supervisord-zotonic.conf
[program:zotonic]
command=/zotonic/bin/zotonic debug
autorestart=true

